# config/supervisor/producer.conf
[program:producer] 
command=/home/ubuntu/code/spark-jobs-automation/api/run-api.sh
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/uvicorn.err.log
stdout_logfile=/var/log/uvicorn.out.log